# roles/packages/darkman/templates/darkman.service.j2
# {{ ansible_managed }}
[Unit]
Description=Framework for dark-mode and light-mode transitions.
Documentation=man:darkman(1)

[Service]
Type=dbus
BusName=nl.whynothugo.darkman
ExecStart=/usr/local/bin/darkman run
Restart=on-failure
TimeoutStopSec=15
Slice=background.slice

[Install]
WantedBy=default.target
