# roles/packages/gcalcli/templates/systemd.service.j2
# {{ ansible_managed }}
[Unit]
Description=Get next meeting from agenda

[Service]
Type=Oneshot
ExecStart=/home/{{ username }}/bin/next-meeting