# roles/packages/desktop/templates/i3lock-sleep.service.j2
#
# {{ ansible_managed }}
#
[Unit]
Description=i3lock on sleep at sleep time
Before=sleep.target

[Service]
Type=forking
User=%I
Environment=DISPLAY=:0
ExecStart=/usr/bin/i3lock -c 000000

[Install]
WantedBy=suspend.target
