# roles/packages/gh/templates/systemd.timer.j2
# {{ ansible_managed }}
[Unit]
Description=Auto-refresh the gh notification counter

[Timer]
OnBootSec=5min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target