# roles/packages/gh/tasks/templates/systemd.service.j2
# {{ ansible_managed }}
[Unit]
Description=Get GH notifications counter

[Service]
Type=oneshot
ExecStart=/home/{{ username }}/bin/gh-notification-counter